# package/debian/DEBIAN/postinst
#!/bin/bash -e
systemctl enable --now haltonika
#systemctl status haltonika
