# package/debian/DEBIAN/prerm
#!/bin/bash -e
systemctl disable --now haltonika
#systemctl status haltonika